# src/chronos/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/1 * * * * root  . /app/.env && /app/scan_products.py >> /var/log/scan_products.log 2>&1

